# package/Linux/gateway/rpm/prerm
#!/bin/sh

set -e

# $1 == 1 during upgrades; $1 == 0 on final erase. Skip shutdown for upgrades.
# Equivalent to %systemd_preun: disable --now on erase only.
if [ "$1" -eq 0 ] && [ -d /run/systemd/system ]; then
	systemctl --no-reload disable --now devolutions-gateway.service >/dev/null 2>&1 || true
fi
